# build/wackamole/wackamole.xml
<?xml version='1.0'?>
<!DOCTYPE service_bundle SYSTEM '/usr/share/lib/xml/dtd/service_bundle.dtd.1'>
<service_bundle type='manifest' name='export'>
  <service name='/omniti/network/wackamole' type='service' version='0'>
    <dependency name='filesystem_minimal' grouping='require_all' restart_on='none' type='service'>
      <service_fmri value='svc:/system/filesystem/minimal'/>
    </dependency>
    <dependency name='loopback' grouping='require_any' restart_on='error' type='service'>
      <service_fmri value='svc:/network/loopback'/>
    </dependency>
    <dependency name='network' grouping='optional_all' restart_on='error' type='service'>
      <service_fmri value='svc:/milestone/network'/>
    </dependency>
    <dependency name='spread' grouping='require_all' restart_on='error' type='service'>
      <service_fmri value='svc:/network/spread'/>
    </dependency>
    <exec_method name='stop' type='method' exec=':kill' timeout_seconds='60'>
      <method_context/>
    </exec_method>
    <property_group name='tm_common_name' type='template'/>
    <instance name='wackamole' enabled='false'>
      <dependency name='config_data' grouping='require_all' restart_on='none' type='path'>
        <service_fmri value='file://localhost/opt/omni/etc/wackamole.conf'/>
      </dependency>
      <exec_method name='start' type='method' exec='/opt/omni/sbin/wackamole' timeout_seconds='60'>
        <method_context>
          <method_credential user='root' group='root' privileges='basic,!proc_session,!proc_info,!file_link_any,sys_net_config,net_icmpaccess,net_rawaccess,net_privaddr,file_dac_read,file_dac_search,sys_resource'/>
        </method_context>
      </exec_method>
      <property_group name='general' type='framework'>
        <propval name='action_authorization' type='astring' value='solaris.smf.manage.bind'/>
        <propval name='value_authorization' type='astring' value='solaris.smf.manage.bind'/>
      </property_group>
    </instance>
    <stability value='Unstable'/>
  </service>
</service_bundle>
