# build/apache22/files/manifest-http-apache.xml
<?xml version='1.0' ?>
<!DOCTYPE service_bundle SYSTEM '/usr/share/lib/xml/dtd/service_bundle.dtd.1'>

<service_bundle type='manifest' name='OMNIapache:httpd' >
	<service name='network/http' type='service' version='1'>
		<instance name='apache22' enabled='false'>

			<dependency name='loopback'  grouping='require_all'
			restart_on='error' type='service'>
				<service_fmri value='svc:/network/loopback:default' />
			</dependency>
			<dependency name='physical'  grouping='optional_all'
			restart_on='error' type='service'>
				<service_fmri value='svc:/network/physical:default' />
			</dependency>

			<dependency name='localfs'  grouping='require_all'
			restart_on='error' type='service'>
				<service_fmri value='svc:/system/filesystem/local:default' />
			</dependency>
			<dependency
			name='config-file'
			grouping='require_all'
			restart_on='none'
			type='path'>
			<service_fmri
			value='file://localhost/opt/apache22/conf/httpd.conf' />
			</dependency>

			<exec_method name='start' type='method'
			exec='/opt/apache22/bin/httpd -f /opt/apache22/conf/httpd.conf -k start'
			timeout_seconds='60' />

			<exec_method name='stop' type='method'
			exec='/opt/apache22/bin/httpd -f /opt/apache22/conf/httpd.conf -k stop'
			timeout_seconds='60' />

			<exec_method name='refresh' type='method'
			exec='/opt/apache22/bin/httpd -f /opt/apache22/conf/httpd.conf -k graceful'
			timeout_seconds='60' />

		<property_group name='general' type='framework'>
			<propval name='action_authorization' type='astring'
			value='omniti.smf.manage.http/apache' />
		</property_group>

	</instance>
<stability value='Evolving' />
<template>
	<common_name>
		<loctext xml:lang='C'>Apache HTTP Server </loctext>
	</common_name>
	<documentation>
		<manpage title='apache' section='1M' />
		<doc_link name='apache.org' uri='http://httpd.apache.org' />
	</documentation>
</template>
</service>
</service_bundle>
