# build/isc-dhcp/files/isc-dhcp.xml
<?xml version='1.0' ?>
<!DOCTYPE service_bundle SYSTEM '/usr/share/lib/xml/dtd/service_bundle.dtd.1'>
<!--

	This file and its contents are supplied under the terms of the
	Common Development and Distribution License ("CDDL"), version 1.0.
	You may only use this file in accordance with the terms of version
	1.0 of the CDDL.

	A full copy of the text of the CDDL should have accompanied this
	source.  A copy of the CDDL is also available via the Internet at
	http://www.illumos.org/license/CDDL.

	Copyright 2014 OmniTI Computer Consulting, Inc. All rights reserved.

-->

<!--	XXX KEBE SAYS we need more properties, and we may need a method
	script to start these running.
-->

<service_bundle type='manifest' name='isc-dhcp'>

<service name='network/service/dhcp' type='service' version='1'>
	<dependency name='loopback'  grouping='require_all'
	restart_on='none' type='service'>
		<service_fmri value='svc:/network/loopback:default' />
	</dependency>

	<dependency name='physical'  grouping='require_all'
	restart_on='none' type='service'>
		<service_fmri value='svc:/network/physical:default' />
	</dependency>

	<dependency name='localfs'  grouping='require_all'
	restart_on='none' type='service'>
		<service_fmri value='svc:/system/filesystem/local:default' />
	</dependency>

	<instance name='ipv4' enabled='false'>
	<exec_method name='start' type='method'
	exec='/usr/sbin/dhcpd -q -cf %{config_file} -lf %{lease_file} -p %{port} -s %{address} %{listen_ifnames}'
	timeout_seconds='60' />
	
	<exec_method name='stop' type='method'
	exec=':kill'
	timeout_seconds='60' />
	
	<exec_method name='refresh' type='method'
	exec=':kill -HUP'
	timeout_seconds='60' />

	<property_group name='general' type='framework'>
		<propval name='action_authorization' type='astring'
		value='omnios.smf.manage.dhcp' />
	</property_group>

        <property_group name='application' type='application'>
		<propval name='port' type='integer' value='67' />
		<propval name='address' type='astring'
		value='255.255.255.255' />
		<!-- For least-suprise from other ISC DHCP deployments. -->
		<!-- OI or Oracle Solaris folks know how to fix these. -->
		<propval name='config_file' type='astring'
		value='/etc/dhcpd.conf' />
		<propval name='lease_file' type='astring'
		value='/var/db/dhcpd.leases' />
		<!-- Slight incompatibility. Use space-separated ifnames
		     in one astring instead of an astring_list.
		-->
		<propval name='listen_ifnames' type='astring' value='' />
	</property_group>

	<template>
		<common_name>
			<loctext xml:lang='C'>ISC DHCP server</loctext>
		</common_name>

		<documentation>
			<manpage title='dhcpd' section='8' />
			<doc_link name='isc.org' uri='http://isc.org' />
		</documentation>
	</template>
	</instance>

	<instance name='ipv6' enabled='false'>
	<exec_method name='start' type='method'
	exec='/usr/sbin/dhcpd -6 -q -cf %{config_file} -lf %{lease_file} -p %{port} %{listen_ifnames}'
	timeout_seconds='60' />
	
	<exec_method name='stop' type='method'
	exec=':kill'
	timeout_seconds='60' />
	
	<exec_method name='refresh' type='method'
	exec=':kill -HUP'
	timeout_seconds='60' />

	<property_group name='general' type='framework'>
		<propval name='action_authorization' type='astring'
		value='omnios.smf.manage.dhcp' />
	</property_group>

        <property_group name='application' type='application'>
		<propval name='port' type='integer' value='67' />
		<!-- Slight incompatibility. Use space-separated ifnames
		     in one astring instead of an astring_list.
		-->
		<propval name='listen_ifnames' type='astring' value="" />
		<propval name='config_file' type='astring'
		value='/etc/dhcpd6.conf' />
		<propval name='lease_file' type='astring'
		value='/var/db/dhcpd6.leases' />
	</property_group>

	<template>
		<common_name>
			<loctext xml:lang='C'>ISC DHCPv6 server</loctext>
		</common_name>

		<documentation>
			<manpage title='dhcpd' section='8' />
			<doc_link name='isc.org' uri='http://isc.org' />
		</documentation>
	</template>
	</instance>

	<stability value='Evolving' />

</service>

<service name='network/service/dhcrelay' type='service' version='1'>
	<dependency name='loopback'  grouping='require_all'
	restart_on='none' type='service'>
		<service_fmri value='svc:/network/loopback:default' />
	</dependency>

	<dependency name='physical'  grouping='require_all'
	restart_on='none' type='service'>
		<service_fmri value='svc:/network/physical:default' />
	</dependency>

	<dependency name='localfs'  grouping='require_all'
	restart_on='none' type='service'>
		<service_fmri value='svc:/system/filesystem/local:default' />
	</dependency>

	<instance name='ipv4' enabled='false'>
	<exec_method name='start' type='method'
	exec='/lib/svc/method/dhcrelay -4 %{hop_count} %{append_agent_option} %{listen_ifnames} %{servers}'
	timeout_seconds='60' />
	
	<exec_method name='stop' type='method'
	exec=':kill'
	timeout_seconds='60' />
	
	<exec_method name='refresh' type='method'
	exec=':kill -HUP'
	timeout_seconds='60' />

	<property_group name='general' type='framework'>
		<propval name='action_authorization' type='astring'
		value='omnios.smf.manage.dhcp' />
	</property_group>

        <property_group name='application' type='application'>
		<propval name='hop_count' type='integer' value='10' />
		<propval name='append_agent_option' type='boolean'
		value='false'/>
		<!-- A space-separated server list is here for now. -->
		<propval name='servers' type='astring' value='' />
		<!-- A space-separated interface list is here for now. -->
		<propval name='listen_ifnames' type='astring' value='' />
        </property_group>

	<template>
		<common_name>
			<loctext xml:lang='C'>ISC DHCP relay</loctext>
		</common_name>

		<documentation>
			<manpage title='dhcrelay' section='8' />
			<doc_link name='isc.org' uri='http://isc.org' />
		</documentation>
	</template>
	</instance>

	<instance name='ipv6' enabled='false'>
	<exec_method name='start' type='method'
	exec='/lib/svc/method/dhcrelay -6 %{hop_count} NOP %{receive_query_links} %{forward_query_links}'
	timeout_seconds='60' />
	
	<exec_method name='stop' type='method'
	exec=':kill'
	timeout_seconds='60' />
	
	<exec_method name='refresh' type='method'
	exec=':kill -HUP'
	timeout_seconds='60' />

	<property_group name='general' type='framework'>
		<propval name='action_authorization' type='astring'
		value='omnios.smf.manage.dhcp' />
	</property_group>

        <property_group name='application' type='application'>
		<propval name='hop_count' type='integer' value='10' />
		<!-- A space-separated receive list is here for now. -->
		<propval name='receive_query_links' type='astring' value='' />
		<!-- A space-separated forward list is here for now. -->
		<propval name='forward_query_links' type='astring' value='' />
        </property_group>

	<template>
		<common_name>
			<loctext xml:lang='C'>ISC DHCPv6 relay</loctext>
		</common_name>

		<documentation>
			<manpage title='dhcpd' section='8' />
			<doc_link name='isc.org' uri='http://isc.org' />
		</documentation>
	</template>
	</instance>

	<stability value='Evolving' />

</service>

</service_bundle>
