# build/openssh/ssh.xml
<?xml version='1.0' ?>
<!DOCTYPE service_bundle SYSTEM '/usr/share/lib/xml/dtd/service_bundle.dtd.1'>

<service_bundle type='manifest' name='openssh' >
	<service name='network/ssh' type='service' version='1'>
		<create_default_instance enabled='false' />
		<single_instance />

		<dependency name='loopback'  grouping='require_all'
		restart_on='none' type='service'>
			<service_fmri value='svc:/network/loopback:default' />
		</dependency>

		<dependency name='physical'  grouping='require_all'
		restart_on='none' type='service'>
			<service_fmri value='svc:/network/physical:default' />
		</dependency>

		<dependency name='cryptosvc'  grouping='require_all'
		restart_on='none' type='service'>
			<service_fmri value='svc:/system/cryptosvc:default' />
		</dependency>
		
		<dependency name='localfs'  grouping='require_all'
		restart_on='none' type='service'>
			<service_fmri value='svc:/system/filesystem/local:default' />
		</dependency>

		<dependency
		name='config-file'
		grouping='require_all'
		restart_on='restart'
		type='path'>
		<service_fmri
		value='file://localhost/etc/ssh/sshd_config' />
		</dependency>
		
		<exec_method name='start' type='method'
		exec='/lib/svc/method/sshd start'
		timeout_seconds='60' />
		
		<exec_method name='stop' type='method'
		exec=':kill'
		timeout_seconds='60' />
		
		<exec_method name='refresh' type='method'
		exec=':kill -HUP'
		timeout_seconds='60' />

		<property_group name='general' type='framework'>
			<propval name='action_authorization' type='astring'
			value='omnios.smf.manage.openssh' />
		</property_group>

		<stability value='Evolving' />

		<template>
			<common_name>
				<loctext xml:lang='C'>OpenSSH Server</loctext>
			</common_name>

			<documentation>
				<manpage title='sshd' section='1M' />
				<doc_link name='openssh.com' uri='http://openssh.com' />
			</documentation>
		</template>
	</service>
</service_bundle>
