# build/rabbitmq/files/rabbitmq.xml
<?xml version='1.0'?>
<!DOCTYPE service_bundle SYSTEM '/usr/share/lib/xml/dtd/service_bundle.dtd.1'>
<service_bundle type='manifest' name='export'>
  <service name='network/rabbitmq' type='service' version='0'>
    <create_default_instance enabled='false'/>
    <single_instance/>
    <dependency name='network' grouping='require_all' restart_on='none' type='service'>
      <service_fmri value='svc:/milestone/network:default'/>
    </dependency>
    <dependency name='config-file' grouping='require_all' restart_on='none' type='path'>
      <service_fmri value='file://localhost/opt/rabbitmq/etc/rabbitmq/rabbitmq-env.conf'/>
    </dependency>
    <method_context working_directory="/opt/rabbitmq/var/lib/rabbitmq">
      <method_environment>
        <envvar name="HOME" value="/opt/rabbitmq/var/lib/rabbitmq"/>
        <envvar name="PATH" value="/usr/bin:/usr/sbin:/opt/omni/bin"/>
      </method_environment>
    </method_context>
    <exec_method name='start' type='method' exec='/opt/rabbitmq/sbin/rabbitmq-server -detached' timeout_seconds='60'>
      <method_context>
        <method_credential user='rabbitmq' group='rabbitmq' />
      </method_context>
    </exec_method>
    <exec_method name='stop' type='method' exec='/opt/rabbitmq/sbin/rabbitmqctl stop' timeout_seconds='60'>
      <method_context>
        <method_credential user='rabbitmq' group='rabbitmq' />
      </method_context>
    </exec_method>
    <stability value='Evolving'/>
    <template>
      <common_name>
        <loctext xml:lang='C'>RabbitMQ</loctext>
      </common_name>
      <documentation>
        <manpage title='rabbitmq' section='8'/>
      </documentation>
    </template>
  </service>
</service_bundle>
