# build/virtualbox/virtualbox-kerneldriver.xml
<?xml version='1.0'?>
<!DOCTYPE service_bundle SYSTEM "/usr/share/lib/xml/dtd/service_bundle.dtd.1">

<service_bundle type='manifest' name='SUNWvbox:kerneldriver'>

<service
    name='application/virtualbox/kerneldriver'
    type='service'
    version='1'>

    <create_default_instance enabled='true' />

    <single_instance />
  
    <!-- Wait for local filesystems to be mounted (just to be safe, don't start too early) -->
    <dependency
        name='filesystem-local'
        grouping='require_all'
        restart_on='none'
        type='service'>
        <service_fmri value='svc:/system/filesystem/local:default' />
    </dependency>

    <exec_method
        type='method'
        name='start'
        exec='/opt/VirtualBoxKernel/vboxconfig.sh --installdrivers'
        timeout_seconds='20' >
        <method_context>
            <method_credential user='root' group='root' />
        </method_context>
    </exec_method>

    <exec_method
        type='method'
        name='stop'
        exec='/opt/VirtualBoxKernel/vboxconfig.sh --removedrivers'
        timeout_seconds='60' >
        <method_context>
            <method_credential user='root' group='root' />
        </method_context>
    </exec_method>

    <property_group name='startd' type='framework'>
      <propval name='duration' type='astring' value='transient' />
    </property_group>

    <template>
        <common_name>
            <loctext xml:lang='C'>VirtualBox Driver Installer</loctext>
        </common_name>
    </template>
</service>

</service_bundle>

